# schedule/insider-backup.timer
[Unit]
Description=Insider Tracker — Daily PostgreSQL backup (05:30 UTC, after nightly ingest + sync_job)

[Timer]
OnCalendar=*-*-* 05:30 UTC
Persistent=true
Unit=insider-backup.service

[Install]
WantedBy=timers.target
