# schedule/insider-backtest.timer
[Unit]
Description=Run insider backtest daily after backup completes

[Timer]
OnCalendar=*-*-* 06:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
