# schedule/insider-perf-profile.timer
[Unit]
Description=Run insider performance profile refresh weekly

[Timer]
OnCalendar=Sun 07:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
